# modules/metacpan/templates/init_watcher.erb
#!/bin/sh
#
### BEGIN INIT INFO
# Provides:          <%= @filename %>
# Required-Start:    $ALL
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: <%= @description %>
### END INIT INFO

# Check for missing binaries (stale symlinks should not happen)
# Note: Special treatment of stop for LSB conformance
export PATH=<%= @perlbin %>:$PATH
WATCHERBIN=<%= @app_root %>/bin/metacpan
PIDFILE=/var/run/<%= @filename %>.pid
ROOT=<%= @app_root %>
USER=metacpan
GROUP=users

DAEMON="<%= app_root %>/bin/carton"
DAEMON_ARGS="exec $WATCHERBIN watcher"

test -x $WATCHERBIN || { echo "$WATCHERBIN not installed";
	if [ "$1" = "stop" ]; then exit 0;
	else exit 5; fi; }

. /lib/lsb/init-functions

case "$1" in
    start)
	echo -n "Starting metacpan-watcher "
	HOME=/home/metacpan /sbin/start-stop-daemon -b -m --start --pidfile $PIDFILE \
	  -c $USER:$GROUP --chdir $ROOT --startas $DAEMON -- $DAEMON_ARGS
	;;
    stop)
	echo -n "Shutting down metacpan-watcher "
    /sbin/start-stop-daemon --stop --oknodo --pidfile $PIDFILE

	;;
    restart)
	## Stop the service and regardless of whether it was
	## running or not, start it again.
	$0 stop
	$0 start

	;;
    status)
	echo -n "Checking for service metacpan-watcher "
    [ -s $PIDFILE ] && kill -0 $(cat $PIDFILE) >/dev/null 2>&1
	;;
    *)
	echo "Usage: $0 {start|stop|status|restart}"
	exit 1
	;;
esac
exit $?
