# modules/starman/templates/init.sh.erb
#!/bin/sh

# Installed by puppet.

# Generated at Thu Jun  5 11:04:48 2014 with Daemon::Control 0.001005

### BEGIN INIT INFO
# Provides:          <%= @service %>
# Required-Start:    $all
# Required-Stop:
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Starts <%= @service %>
# Description:       Starts <%= @service %>
### END INIT INFO`

if [ -x <%= @script %> ];
then
    <%= @script %> $1
else
    echo "Required program <%= @script %> not found!"
    exit 1;
fi
